# crates/gateway/src/service/systemd.service
[Unit]
Description=Crabtalk Gateway - {description}
After=network-online.target

[Service]
Type=simple
WorkingDirectory={config_dir}
Environment=PATH={path}
ExecStart={binary} {subcommand} run --daemon {socket} --config {config_path}
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
